# chroota.zdb/host-initd.zdb/main.ini
########################### host-initd-info
[info]
apply=Install /etc/init.d/{{global_name}} script and register it with update-rc.d
destroy=Remove /etc/init.d/{{global_name}} script and deregister it with update-rc.d

########################### host-initd-script
[file]
path=/etc/init.d/{{global_name}}
mode=+x
content="
#! /bin/sh
# auto generated by zapusk-tool chroota

### BEGIN INIT INFO
# Provides:             {{global_name}}
# Required-Start:       $all
# Required-Stop:        $remote_fs $syslog
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    {{global_name}} chruta
### END INIT INFO

set -e

uninstall() {
  echo -n "Are you really sure you want to uninstall this service script? (call stop or off before that). [yes|No] "
  local SURE
  read SURE
  if [ "$SURE" = "yes" ]; then
    # stop
    # todo - надо ли вообще стоп тут делать? а хотя бы off? жизнь покажет..
    echo "Notice: zapusk component still has to be stopped and removed: '{{state_dir}}'" >&2
    update-rc.d -f {{global_name}} remove
    rm -fv "$0"
  fi
}

case "$1" in
  'on' | 'off')
    {{chroot_tool}} "$1" -d {{chroot_dir}}
    ;;

  'chroot') {{chroot_tool}} on -d {{chroot_dir}} && chroot {{chroot_dir}}
    ;;
  'uninstall')
    uninstall
    ;;
  *)
  {{chroot_tool}} run "/zapusk-tool/zapusk \"$1\" --zdb /chroot.d --state_dir /chroot.d/_state --deferred-master" -d "{{chroot_dir}}"
esac
"

############################ host-initd-register
[os]
type=os
apply=update-rc.d {{global_name}} defaults
destroy=update-rc.d -f {{global_name}} remove
#  && rm -fv /etc/init.d/{{global_name}} -- вроде как не надо файл удалять, [file] сам удалит
