# chroota.zdb/host-initd.zdb/main.ini
############################ main #######################
# приоритет 09 перед main, чтобы в хост чрута закреплялась
# независимо от успешности restart или start на этапе установки

[commands]
host
system-update
apply

[once]
apply
dir={{chroot_dir}}

[info]
apply=Install /etc/init.d/{{global_name}} script and register it with update-rc.d
destroy=Remove /etc/init.d/{{global_name}} script and deregister it with update-rc.d

[file]
path=/etc/init.d/{{global_name}}
mode=+x
content="
#! /bin/sh
# auto generated by zapusk-tool chroota

### BEGIN INIT INFO
# Provides:             {{global_name}}
# Required-Start:       $all
# Required-Stop:        $remote_fs $syslog
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    {{global_name}} chruta
### END INIT INFO

set -e

case "$1" in
  'chroot') {{chroot_tool}} on -d {{chroot_dir}} && chroot {{chroot_dir}}
    ;;
  *)
  {{chroot_tool}} run "/zapusk-tool/zapusk \"$1\" --zdb /chroot.d --state_dir /chroot.d/_state" -d "{{chroot_dir}}"
esac
"

# todo [once]

[]
type=os
apply=update-rc.d {{global_name}} defaults
destroy=update-rc.d -f {{global_name}} remove

#  && rm -fv /etc/init.d/{{global_name}} -- вроде как не надо файл удалять, [file] сам удалит
